# packaging/systemd/quickshell-polkit-agent-debug.service
[Unit]
Description=Quickshell Polkit Authentication Agent (Debug Mode)
Documentation=https://github.com/bennyp/quickshell-polkit-agent
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY
# Only run if quickshell is available
ConditionPathExists=/usr/bin/quickshell

[Service]
Type=exec
ExecStart=/usr/libexec/quickshell-polkit-agent
ExecStopPost=/bin/rm -f /run/user/%%i/quickshell-polkit*
Slice=session.slice
TimeoutStartSec=10sec
TimeoutStopSec=5sec
Restart=on-failure
RestartSec=2sec
StartLimitBurst=5
StartLimitIntervalSec=60sec

# Debug environment
Environment=QT_LOGGING_RULES="polkit.*=true;ipc.*=true"
Environment=QT_FORCE_STDERR_LOGGING=1

# Automatic runtime directory cleanup
RuntimeDirectory=quickshell-polkit
RuntimeDirectoryMode=0700

# Minimal security (for debugging)
PrivateTmp=yes
NoNewPrivileges=yes

# Resource limits (more generous for debugging)
MemoryMax=128M
TasksMax=100

[Install]
WantedBy=graphical-session.target