# packaging/systemd/quickshell-polkit-agent.service
[Unit]
Description=Quickshell Polkit Authentication Agent
Documentation=https://github.com/bennyp/quickshell-polkit-agent
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY
# Only run if quickshell is available
ConditionPathExists=/usr/bin/quickshell

[Service]
Type=exec
ExecStart=/usr/libexec/quickshell-polkit-agent
ExecStopPost=/bin/rm -f ${RUNTIME_DIRECTORY}/quickshell-polkit*
Slice=session.slice
TimeoutStartSec=10sec
TimeoutStopSec=5sec
Restart=on-failure
RestartSec=2sec
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
RuntimeDirectory=quickshell-polkit
RuntimeDirectoryMode=0700

# Resource limits
MemoryMax=64M
TasksMax=50

[Install]
WantedBy=graphical-session.target